# templates/webdriver-wharf.service.erb
[Unit]
Description=WebDriver Wharf
Requires=docker

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/webdriver-wharf
ExecStart=/usr/bin/webdriver-wharf

[Install]
WantedBy=multi-user.target
